# roles/restic/files/restic-snapshot.timer
[Timer]
OnBootSec=15min
OnUnitInactiveSec=15min
