# roles/hydroxide/templates/hydroxide.service
[Service]
ExecStart=hydroxide -smtp-host={{ wgme.ip }} -imap-host={{ wgme.ip }} -carddav-host={{ wgme.ip }} serve